# scripts/skeleton-home-edge-executor.service
[Unit]
Description=Skeleton Home Edge executor is intentionally one-shot over strict OpenSSH
Documentation=docs/HOME_EDGE_EXECUTOR.md

# This file is intentionally not an installable long-running service. The
# supported execution route is:
#
#   ssh -o BatchMode=yes -o StrictHostKeyChecking=yes ... home_edge_exec --server
#
# `home_edge_exec --server` reads exactly one signed JSON request from stdin,
# executes it, writes one JSON receipt to stdout, and exits. Running that
# one-shot stdin protocol under a simple restartable service would create
# a broken restart loop with no request input.
